# ops/systemd/traffic-api.service.d/30-memory-guard.conf
[Service]
Restart=always
RestartSec=5
MemoryHigh=2200M
MemoryMax=3000M
