# ops/systemd/traffic-api.service.d/40-exec-python-module.conf
[Service]
ExecStart=
ExecStart=/var/www/Traffic/traffic-api/.venv/bin/python3 -m uvicorn app.main:app --host 127.0.0.1 --port 3345 --timeout-keep-alive 5 --timeout-graceful-shutdown 15
