# ops/systemd/traffic-api.service.d/50-venv-guard.conf
[Service]
ExecStartPre=/var/www/Traffic/ops/traffic-api-venv-guard.sh
