# ops/systemd/traffic-web.service.d/30-restart-memory-guard.conf
[Service]
Restart=always
RestartSec=5
MemoryHigh=700M
MemoryMax=1100M
TimeoutStopSec=20
KillMode=mixed
